# systemd/usage_report.service
[Unit]
Description=Grit Usage Report

[Service]
ExecStart=/usr/local/bin/usage_report
Type=oneshot

[Install]
WantedBy=default.target
